# bin/oracle_agent_stopstart.redhat
#!/bin/bash
#
# oracle_agent_stopstart.redhat          Start Stop oracle agent.
#
# chkconfig: 235 40 60
# description: This script starts and stops the oracle agent.
#
# Works for: RHEL
#
# Activate using chkconfig:
# (as root)
# cp /path/to/tools/bin/oracle_agent_stopstart.redhat /etc/init.d
# chkconfig --add oracle_agent_stopstart.redhat
#
# Frits Hoogland
# -version control
# 16042009 FH Creation
########################################
# Please set the following variables ! #
########################################
#
# TOOLS_HOME
TOOLS_HOME=/oracle/tools
#
# ORACLE_HOME - oracle home of the agent
ORACLE_HOME=/oracle/gc/agent10g
#
# ORACLE_USER - userid with which the agent is installed
ORACLE_USER=oracle
#
# TZ - timezone is import for agent; only reset this if it gives troubles
TZ=$TZ
#
#############################################
# Do not modify anything beyond this line ! #
#############################################

# Source function library.
. /etc/init.d/functions

start() {
        touch /var/lock/subsys/oracle_agent_stopstart.redhat
        cat $TOOLS_HOME/etc/agenttab | while read LINE; do
                case $LINE in
                        \#*)
                                ;;
                        *)
                                if [ "$( echo $LINE | awk -F: '{print $3 }' )" = "Y" ]; then
                                        AGENT_NAME=$( echo $LINE | awk -F: '{print $1}' )
                                        ORACLE_HOME=$( echo $LINE | awk -F: '{print $2}' )
                                        echo -n $"Starting oracle agent $AGENT_NAME: "
                                        su $ORACLE_USER -c "export TZ=$TZ; export ORACLE_HOME=$ORACLE_HOME; $ORACLE_HOME/bin/emctl start agent 2>&1 > /dev/null"
					RETVAL=$?
                                        [ "$RETVAL" -gt 0 ] && echo $"FAILED" || echo $"OK"
					return $RETVAL
                                fi
                                ;;
                esac
        done
}
stop() {
        rm -f /var/lock/subsys/oracle_agent_stopstart.redhat
        cat $TOOLS_HOME/etc/agenttab | while read LINE; do
                case $LINE in
                        \#*)
                                ;;
                        *)
                                if [ "$( echo $LINE | awk -F: '{print $3 }' )" = "Y" ]; then
                                        AGENT_NAME=$( echo $LINE | awk -F: '{print $1}' )
                                        ORACLE_HOME=$( echo $LINE | awk -F: '{print $2}' )
                                        echo -n $"Stopping oracle agent $AGENT_NAME: "
                                        su $ORACLE_USER -c "export TZ=$TZ; export ORACLE_HOME=$ORACLE_HOME; $ORACLE_HOME/bin/emctl stop agent 2>&1 > /dev/null"
					RETVAL=$?
                                        [ "$RETVAL" -gt 0 ] && echo $"FAILED" || echo $"OK"
					return $RETVAL
                                fi
                                ;;
                esac
        done
}
status() {
        cat $TOOLS_HOME/etc/agenttab | while read LINE; do
                case $LINE in
                        \#*)
                                ;;
                        *)
                                if [ "$( echo $LINE | awk -F: '{print $3 }' )" = "Y" ]; then
                                        AGENT_NAME=$( echo $LINE | awk -F: '{print $1}' )
                                        ORACLE_HOME=$( echo $LINE | awk -F: '{print $2}' )
                                        echo $"Status of oracle agent $AGENT_NAME: "
                                        su $ORACLE_USER -c "export ORACLE_HOME=$ORACLE_HOME; $ORACLE_HOME/bin/emctl status agent"
                                fi
                                ;;
                esac
        done
}
restart() {
        stop
        start
}


case "$1" in
  start)
        start
        ;;
  stop)
        stop
        ;;
  restart)
        restart
        ;;
  status)
        status
        ;;
  *)
        echo $"Usage: $0 {start|stop|restart|status}"
        exit 1
        ;;
esac
