# etc/cron.daily/oracle_asm_audit_tmpwatch
#!/bin/bash
# 168 = 24 * 7; a week
HOURS=168
for DIR in /u01/app/*/grid/rdbms/audit; do
  if [ -d "$DIR" ]; then
    /usr/sbin/tmpwatch -maf $HOURS "$DIR"
  fi
done
